# packaging/arch/mousetrap-hyprland.install
post_install() {
  cat <<'EOF'
Mousetrap Hyprland installed.

Next steps:
  1. Enable click backend:
       systemctl --user enable --now ydotool.service

  2. Create default config:
       mousetrap-hyprland init-config

  3. Add a Hyprland bind, for example:
       bind = SUPER, SPACE, exec, mousetrap-hyprland-activate

     Or copy the example config from:
       /usr/share/doc/mousetrap-hyprland/mousetrap.conf
EOF
}

post_upgrade() {
  post_install
}
